# redi301.conf
# supervisor configuration file
[program:redi301]
process_name=redi301
directory=/root/redi301
command=/root/redi301/redi301 -a :80
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/redi301.log